# automation/research/claim_verifier/claim-process.service
[Unit]
Description=AIで投資の壁を越える 主張抽出 + T=0 スナップショット
After=network-online.target research-collector.service

[Service]
Type=oneshot
WorkingDirectory=/home/sol/aiquant-lab
ExecStart=/home/sol/aiquant-lab/automation/research/claim_verifier/run_process.sh
Environment=HOME=/home/sol
Environment=PATH=/home/sol/.local/bin:/home/sol/.claude/local:/usr/local/bin:/usr/bin:/bin
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1200
