# automation/research/claim_verifier/claim-verify.service
[Unit]
Description=AIで投資の壁を越える 期限切れ主張の検証

[Service]
Type=oneshot
WorkingDirectory=/home/sol/aiquant-lab
ExecStart=/home/sol/aiquant-lab/automation/research/claim_verifier/run_verify.sh
Environment=HOME=/home/sol
Environment=PATH=/home/sol/.local/bin:/home/sol/.claude/local:/usr/local/bin:/usr/bin:/bin
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
