# automation/research/claim_verifier/claim-report.timer
[Unit]
Description=主張検証 週次レポート (毎週月曜 06:00)

[Timer]
OnCalendar=Mon *-*-* 06:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
