# automation/research/research-collector.service
[Unit]
Description=AIで投資の壁を越える 日次リサーチコレクター
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/sol/aiquant-lab
ExecStart=/home/sol/aiquant-lab/automation/research/run.sh
Environment=HOME=/home/sol
Environment=PATH=/home/sol/.local/bin:/home/sol/.claude/local:/usr/local/bin:/usr/bin:/bin
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600
